# cd-rip-eject.service
# NOTE: ExecStart path is updated by Install-cd-rip.sh
# to match the actual installation directory.

[Unit]
Description=Audio cd ejection
After=udisks-glue

[Service]
Type=simple

RemainAfterExit=no

# This path is updated automatically by Install-cd-rip.sh
ExecStart=__INSTALL_DIR__/cd-rip-eject.sh

ExecStop=killall -q cd-rip-and-or-play.sh

# StandardOutput=null

# StandardError=syslog

[Install]
WantedBy=multi-user.target
