# cd-rip-and-or-play.service.in
# We will use a oneshot service that can be called but doesn't persist. This is
# needed because tasks called directly from udev that take longer than 60s get timed out.
#
# NOTE: WorkingDirectory and ExecStart paths are updated by Install-cd-rip.sh
# to match the actual installation directory.

[Unit]
Description=CD rip and-or play on insertion
After=udisks-glue

[Service]
Type=oneshot

RemainAfterExit=no

# These paths are updated automatically by Install-cd-rip.sh
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/cd-rip-and-or-play.sh

ExecStop=killall -q abcde

# StandardOutput=null

# StandardError=syslog

[Install]
WantedBy=multi-user.target
